# ansible/roles/common/templates/systemd_simple.service.j2
[Unit]
Description={{ service_description }}
After={{ service_after }}

[Service]
ExecStart={{ exec_start }}


[Install]
WantedBy=default.target
